# data/systemd/fan-pid.service
[Unit]
Description=PID-based EMC2305 fan controller for LS1046A Mono Gateway DK
Documentation=file:///usr/local/bin/fan-pid
DefaultDependencies=no
# Run after sysfs/hwmon devices exist but before normal multi-user services.
After=sysinit.target systemd-modules-load.service
Before=multi-user.target shutdown.target
Conflicts=shutdown.target fancontrol.service
# Board gate: only Mono Gateway DK has the EMC2305 + LS1046A thermal zones.
ConditionPathExists=/proc/device-tree/compatible

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-pid
Restart=on-failure
RestartSec=5
# Daemon already handles SIGTERM cleanly (writes PWM=MAX as failsafe).
KillSignal=SIGTERM
TimeoutStopSec=10
# Daemon logs via syslog(); journald collects it under SYSLOG_IDENTIFIER=fan-pid.
StandardOutput=journal
StandardError=journal
# Hardening: needs RW on /sys/class/hwmon/.../pwm1 only; no network, no exec.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
NoNewPrivileges=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# /sys is mounted ro under ProtectSystem=strict; hwmon writes need explicit rw.
ReadWritePaths=/sys/class/hwmon

[Install]
WantedBy=multi-user.target