# data/systemd/fan-pid.service
[Unit]
Description=PID-based EMC2305 fan controller for LS1046A Mono Gateway DK
Documentation=file:///usr/local/bin/fan-pid
# Run after sysfs/hwmon devices exist; we keep DefaultDependencies=yes so
# multi-user.target's WantedBy= symlink actually pulls us in at boot.
# (DefaultDependencies=no detached the unit from the target graph and the
#  .wants/ symlink alone was not enough — observed 2026-05-11.)
After=sysinit.target systemd-modules-load.service
Conflicts=fancontrol.service
# Board gate: only Mono Gateway DK has the EMC2305 + LS1046A thermal zones.
# Match the LS1046A SoC string (devicetree compatible is a NUL-separated list,
# ConditionPathExistsGlob can't peek inside; the path-exists check is gated by
# /proc/device-tree being a real DT node, which only LS1046A-class boards expose).
ConditionPathExists=/proc/device-tree/compatible
ConditionPathExistsGlob=/sys/bus/i2c/drivers/emc2305/*-002e

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-pid
Restart=on-failure
RestartSec=5
# Daemon already handles SIGTERM cleanly (writes PWM=MAX as failsafe).
KillSignal=SIGTERM
TimeoutStopSec=10
# Daemon logs via syslog(); journald collects it under SYSLOG_IDENTIFIER=fan-pid.
StandardOutput=journal
StandardError=journal
# Hardening: needs RW on /sys/class/hwmon/.../pwm1 only; no network, no exec.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
NoNewPrivileges=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# /sys is mounted ro under ProtectSystem=strict; hwmon writes need explicit rw.
ReadWritePaths=/sys/class/hwmon

[Install]
WantedBy=multi-user.target