# board/systemd/fan-pid.service
[Unit]
Description=PID-based EMC2305 fan controller for LS1046A Mono Gateway DK
Documentation=file:///usr/local/bin/fan-pid
# Run after sysfs/hwmon devices exist; we keep DefaultDependencies=yes so
# multi-user.target's WantedBy= symlink actually pulls us in at boot.
# (DefaultDependencies=no detached the unit from the target graph and the
#  .wants/ symlink alone was not enough — observed 2026-05-11.)
After=sysinit.target systemd-modules-load.service
Conflicts=fancontrol.service
# Board gate: only Mono Gateway DK has the EMC2305 + LS1046A thermal zones.
# We gate ONLY on /proc/device-tree/compatible (a stable DT-only path that
# only LS1046A-class boards expose). The previous gate also required
# `ConditionPathExistsGlob=/sys/bus/i2c/drivers/emc2305/*-002e` — that path
# is created asynchronously when the i2c-7 bus probe binds the emc2305
# driver, which RACES with multi-user.target. When the condition fails,
# systemd silently drops the unit and does NOT re-evaluate it later, so
# the daemon would never start even though the device appears moments
# afterward. Observed 2026-05-11: condition met at login time but boot
# journal had `-- No entries --` for fan-pid.service. The hwmon-device
# discovery is now retried inside the daemon at startup with a short
# bounded wait, and a udev rule (10-emc2305-fan-pid.rules) ALSO `systemctl
# start`s us when the device binds — so even if multi-user.target races
# us, the udev path catches us. The DT board gate alone is enough to keep
# us off non-LS1046A hardware.
ConditionPathExists=/proc/device-tree/compatible

[Service]
Type=simple
ExecStart=/usr/local/bin/fan-pid
Restart=on-failure
RestartSec=5
# Daemon already handles SIGTERM cleanly (writes PWM=MAX as failsafe).
KillSignal=SIGTERM
TimeoutStopSec=10
# Daemon logs via syslog(); journald collects it under SYSLOG_IDENTIFIER=fan-pid.
StandardOutput=journal
StandardError=journal
# Hardening: needs RW on /sys/class/hwmon/.../pwm1 only; no network, no exec.
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateNetwork=yes
NoNewPrivileges=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
# /sys is mounted ro under ProtectSystem=strict; hwmon writes need explicit rw.
ReadWritePaths=/sys/class/hwmon

[Install]
WantedBy=multi-user.target