# board/systemd/serial-console-debug.service
# serial-console-debug.service — one-shot diagnostic for serial-getty@ttyS0
#
# Captures effective unit, drop-ins, agetty process state, /dev/ttyS0
# owners, termios, journal tail, and dmesg tail at +20 s after
# multi-user.target. Everything goes to the journal under this unit
# name; retrieve via:
#
#     journalctl -u serial-console-debug --no-pager
#
# Designed so an operator who CAN reach the box via SSH can collect
# all the data we need to diagnose serial-console failures even when
# the console itself is dead.

[Unit]
Description=One-shot serial-getty@ttyS0 diagnostic dump
After=multi-user.target
DefaultDependencies=no
ConditionPathExists=/usr/local/bin/serial-console-debug.sh

[Service]
Type=oneshot
RemainAfterExit=yes
# Run after the late-boot printk burst would have settled — the very
# class of event we're trying to diagnose.
ExecStartPre=/bin/sleep 20
ExecStart=/usr/local/bin/serial-console-debug.sh
StandardOutput=journal
StandardError=journal
# Never fail — this is diagnostic, not load-bearing.
SuccessExitStatus=0 1 2

[Install]
WantedBy=multi-user.target